# etc/init/su2.conf
description "Su2"

start on (filesystem)
stop on runlevel [016]

respawn
setuid data
setgid data
chdir /home/data/openhds-su2

exec /home/data/openhds-su2/env/bin/gunicorn -b 0.0.0.0:4000 su2:app
